# test/integration/templates/supervisor.site.conf
[program:runs]
command=$JAWS_BIN_DIR/run-daemon
numprocs=1
startsecs=3
startretries=10
stopasgroup=true
autorestart=unexpected

[program:transfers]
command=$JAWS_BIN_DIR/transfer-daemon
numprocs=1
startsecs=3
startretries=10
stopasgroup=true
autorestart=unexpected

[program:perf-metrics]
command=$JAWS_BIN_DIR/perf-metrics-daemon
numprocs=1
startsecs=3
startretries=10
stopasgroup=true
autorestart=unexpected

[program:rpc-server]
command=$JAWS_BIN_DIR/rpc-server
numprocs=1
startsecs=3
startretries=10
stopasgroup=true
autorestart=unexpected

[program:message-consumer]
command=$JAWS_BIN_DIR/message-consumer
numprocs=1
startsecs=3
startretries=10
stopasgroup=true
autorestart=unexpected

[group:jaws-site]
programs=runs,transfers,perf-metrics,rpc-server,message-consumer
priority=999
